# crawler/crawler.conf
[program:master_crawler]
command=python -u master_crawler.py
directory=/home/crawler
stdout_logfile=/home/log
stderr_logfile=/home/log